# worker/supervisor/conf.d/celery.conf
; ================================
;  celery supervisor
; ================================

[program:celery]
command=celery -A tasks worker --loglevel=info --concurrency=1
directory=/home/app/celery

user=root
;numprocs=1
stdout_logfile=/home/app/logs/celery.out
stderr_logfile=/home/app/logs/celery.err
autostart=true
autorestart=true

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true
